# roles/monitoring/templates/nginx_exporter.socket.j2
[Unit]
Description=NGINX Prometheus Exporter Socket

[Socket]
ListenStream=9113

[Install]
WantedBy=sockets.target
